# deploy/systemd/relay.service
[Unit]
Description=Relay agency-pipeline http server (board/api/email.naples.agency)
After=network-online.target docker.service
Wants=network-online.target

[Service]
WorkingDirectory=/root/agency-pipeline
EnvironmentFile=/root/agency-pipeline/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/tsx src/server.ts
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
