# deploy/systemd/relay-digest.timer
[Unit]
Description=Relay daily brief at 07:30 UTC

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
